# systemd/http-org-webfreak-omaps-stats.socket
[Unit]
Description = omaps cdn statistics HTTP server

[Socket]
ListenStream = 2142
Accept = yes

[Install]
WantedBy = sockets.target

